# backend/deploy/systemd/jobpilot-api.service
[Unit]
Description=JobPilot FastAPI Service
After=network.target

[Service]
User=root
WorkingDirectory=/opt/jobpilot/backend
EnvironmentFile=/opt/jobpilot/backend/.env
ExecStart=/opt/jobpilot/.venv/bin/gunicorn -k uvicorn.workers.UvicornWorker app.main:app --bind 0.0.0.0:8000 --workers 2
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
